# deploy/agent-me-dashboard.service
[Unit]
Description=agent-me dashboard (Phase 4 — read-only web view)
Documentation=https://github.com/thanhpt1110/agent-me/blob/main/design/dashboard-design.md
After=network-online.target agent-me-bridge.service
Wants=network-online.target
# We don't `Requires=agent-me-bridge.service` on purpose: the dashboard is
# read-only and degrades gracefully if the bridge is offline (it just
# shows zero stats). Keeping units independent means a dashboard restart
# doesn't disturb the bridge.

# This is a USER unit (systemd --user). Install via scripts/install-dashboard.sh.
# Same rationale as agent-me-bridge.service for picking user-level over
# system-level — see that file's preamble.

[Service]
Type=simple
WorkingDirectory=%h/agent-me
EnvironmentFile=%h/agent-me/configs/.env

# Bind to 0.0.0.0:8765 — exposed to the NVIDIA-internal reverse proxy
# (`https://agent-me.nvidia.com`, VPN-gated; see
# `design/reverse-proxy-config.md`). The proxy terminates TLS and forwards
# X-Forwarded-* headers; we accept them via `--forwarded-allow-ips=*`
# because the network itself is the trust boundary (only NVIDIA VPN
# clients can reach the proxy in the first place).
#
# Auth model with VPN-gated access: set DASHBOARD_TRUST_NETWORK=1 in
# configs/.env. The entry point will then permit a non-loopback bind
# without DASHBOARD_TOKEN. Optional defense-in-depth: also set
# DASHBOARD_TOKEN to require a per-operator bearer / cookie on top of
# VPN; the middleware handles both transparently.
#
# To run loopback-only (e.g. dev / SSH-tunnel access), drop --host and
# --forwarded-allow-ips to get the safer default.
ExecStart=%h/.local/bin/uv run agent-me-dashboard --host 0.0.0.0 --port 8765 --forwarded-allow-ips=*

# Read-only over the bridge's SQLite, no writes to bridge state. A
# graceful stop just kills the uvicorn worker; nothing to flush.
TimeoutStopSec=10
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-me-dashboard

[Install]
WantedBy=default.target
